# opencode/.config/systemd/user/opencode-automation-healthcheck.service
[Unit]
Description=OpenCode automation startup healthcheck
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/moo/.config/opencode/bin/run-automation-healthcheck
WorkingDirectory=/home/moo

[Install]
WantedBy=default.target
